# Pacemaker/pcs_cluster_36_nodes_feasibility/scripts/s3server@.service
[Unit]
Description=Motr-S3 Core instance Server

[Service]
Type=simple
ExecStart=/root/s3server-service %i
TimeoutStopSec=7sec

[Install]
WantedBy=multi-user.target


